# nhid.service
[Unit]
Description=nhi daemon

[Service]
User=root
LimitMEMLOCK=10M
Restart=always
Type=simple
ExecStart=nhid

[Install]
WantedBy=multi-user.target
